# supervisor.example.ini
[program:tasks]
command = {PATH_TO_VIRTUAL_ENV}/bin/celery worker -A dystwitic.work -E -l INFO -n dystwitic_tasks
environment = PYTHONPATH="{PATH_TO_PROJECT}:$PYTHONPATH"
directory = {PATH_TO_PROJECT}
user = njmattes
stdout_logfile = /tmp/%(program_name)s.log
redirect_stderr = true
stderr_logfile = /tmp/%(program_name)s.err
stopsignal = TERM
autorestart = true

[program:beat]
command = {PATH_TO_VIRTUAL_ENV}/bin/celery beat -A dystwitic.work -l WARN
environment = PYTHONPATH="{PATH_TO_PROJECT}:$PYTHONPATH"
directory = {PATH_TO_PROJECT}
user = njmattes
stdout_logfile = /tmp/%(program_name)s.log
redirect_stderr = true
stderr_logfile = /tmp/%(program_name)s.err
stopsignal = KILL
autorestart = true

[program:stream]
command = {PATH_TO_VIRTUAL_ENV}/bin/python -m dystwitic.work.twitter_stream
environment = PYTHONPATH="/var/www/dystwitic_prod:$PYTHONPATH"
directory = {PATH_TO_PROJECT}
user = njmattes
stdout_logfile = /tmp/%(program_name)s.log
redirect_stderr = true
stderr_logfile = /tmp/%(program_name)s.err
stopsignal = TERM
autorestart = true

[program:web]
command = {PATH_TO_VIRTUAL_ENV}/bin/python -m run
environment = PYTHONPATH="{PATH_TO_PROJECT}:$PYTHONPATH"
directory = {PATH_TO_PROJECT}
user = njmattes
stdout_logfile = /tmp/%(program_name)s.log
redirect_stderr = true
stderr_logfile = /tmp/%(program_name)s.err
stopsignal = KILL
stopasgroup = true
killasgroup = true
autorestart = true

[group:dystwitic]
programs = tasks, beat, stream, web