# supervisor.example.ini
[program:dystwitic_tasks]
command = {PATH_TO_VIRTUAL_ENV}/bin/celery worker -A dystwitic.work -E -l INFO -n dystwitic_tasks
environment = PYTHONPATH="{PATH_TO_PROJECT}:$PYTHONPATH"
directory = {PATH_TO_PROJECT}
user = njmattes
stdout_logfile = /tmp/%(program_name)s.log
redirect_stderr = true
stderr_logfile = /tmp/%(program_name)s.err
stopsignal = TERM
autorestart = true

[program:dystwitic_beat]
command = {PATH_TO_VIRTUAL_ENV}/bin/celery beat -A dystwitic.work -l INFO
environment = PYTHONPATH="{PATH_TO_PROJECT}:$PYTHONPATH"
directory = {PATH_TO_PROJECT}
user = njmattes
stdout_logfile = /tmp/%(program_name)s.log
redirect_stderr = true
stderr_logfile = /tmp/%(program_name)s.err
stopsignal = KILL
autorestart = true

[program:dystwitic_stream]
command = {PATH_TO_VIRTUAL_ENV}/bin/python -m dystwitic.work.twitter_stream
environment = PYTHONPATH="/var/www/dystwitic_prod:$PYTHONPATH"
directory = {PATH_TO_PROJECT}
user = njmattes
stdout_logfile = /tmp/%(program_name)s.log
redirect_stderr = true
stderr_logfile = /tmp/%(program_name)s.err
stopsignal = TERM
autorestart = true

[program:dystwitic_web]
command = {PATH_TO_VIRTUAL_ENV}/bin/python -m run
environment = PYTHONPATH="{PATH_TO_PROJECT}:$PYTHONPATH"
directory = {PATH_TO_PROJECT}
user = njmattes
stdout_logfile = /tmp/%(program_name)s.log
redirect_stderr = true
stderr_logfile = /tmp/%(program_name)s.err
stopsignal = KILL
stopasgroup = true
killasgroup = true
autorestart = true
